# sshd-alert/etc/cron.weekly/alery_cron
#!/bin/bash

# Run the monitoring script weekly
* * * * 0 python3 /usr/local/bin/ssh-login
